# deploy/nearbynoise-web.service
[Unit]
Description=nearbynoise web UI
After=network.target

[Service]
User=www-data
WorkingDirectory=/opt/nearbynoise
# notes_path must be OUTSIDE the rsync archive: the Pi's `rsync -a` resets the
# archive dir's owner/perms every cycle, so www-data cannot write notes there.
ExecStart=/opt/nearbynoise/.venv/bin/python -c "from waitress import serve; from nearbynoise.webserver import create_app; serve(create_app('/var/www/laermprotokoll', '/var/www/laermprotokoll/events.jsonl', notes_path='/var/lib/nearbynoise/notes.json'), host='127.0.0.1', port=8080)"
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
